# services/CallNetworks.service
[Translation]
File=QtopiaServices
Context=CallNetworks
[Service]
Icon = phonenetworks/phonenetworks
Name[]=Call Networks
