# services/Camera.service
[Translation]
File=QtopiaServices
Context=Camera
[Service]
Actions = "getImage(QDSActionRequest);shutter()"
Icon = Camera/camera
Name[]=Camera
