# services/Clock.service
[Translation]
File=QtopiaServices
Context=Clock
[Service]
Actions = "showClock()"
Icon = Clock
Name[]=Show Clock
