# services/PhotoEdit.service
[Translation]
File=QtopiaServices
Context=Photo Edit
[Service]
Actions = "editImage(QDSActionRequest)"
Icon=PhotoEdit
Name[]=Photo Edit
