# services/Ringtone.service
[Translation]
File=QtopiaServices
Context=Ringtone
[Service]
Actions="startMessageRingtone();stopMessageRingtone();startRingtone(QString);stopRingtone(QString)"
Name[]=Ringtone
[startMessageRingtone()]
[stopMessageRingtone()]
[startRingtone(QString)]
[stopRingtone(QString)]
