# services/SettingsManager.service
[Translation]
File=QtopiaServices
Context=CaptureSetting
[Service]
