# services/TimeUpdate.service
[Service]
Actions="updateFromExternalSources(bool,bool,bool,bool)"
