# services/WebAccess.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Actions="openURL(QString)"
